# vagrant-docker-python/docker-image/rootfs/etc/services.d/sshd/run
#!/bin/sh

addgroup -S sshd
adduser -S sshd -G sshd

SSHD=/usr/sbin/sshd

echo "Starting sshd:"
$SSHD -D